# config/supervisor.d/nginx.ini.tmpl
[program:nginx]

{{- if ne .Env.NGINX_DEBUG_ENABLED "false" }}
command=nginx-debug -p /app/var/www -e /dev/stdout -c /opt/etc/nginx/nginx.conf -g 'daemon off;'
{{- else }}
command=nginx -p /app/var/www -e /dev/stdout -c /opt/etc/nginx/nginx.conf -g 'daemon off;'
{{- end }}

stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=false
startretries=0

; Graceful shutdown: SIGQUIT asks nginx for a graceful stop; signal the whole
; process group so workers exit and the listen socket is released cleanly, and
; SIGKILL the group if it overruns stopwaitsecs.
stopsignal=QUIT
stopwaitsecs=20
stopasgroup=true
killasgroup=true
